# pkgcaching/pkgcaching.service
[Unit]
Description=Perform package cache syncing

[Service]
Type=oneshot
ExecStart=/etc/sysconf/caching/run.sh
