# pkgcaching/pkgcaching.service
[Unit]
Description=Perform package cache syncing

[Service]
Type=oneshot
ExecStart=/opt/epiphyte/servers/pkgcache.sh
