# pkgcaching/pkgcaching-server.service
[Unit]
Description=Perform package cache cleanup

[Service]
Type=oneshot
ExecStart=/opt/epiphyte/servers/pkgcache-server.sh
