# pkgcaching/pkgcaching-server.timer
[Unit]
Description=Cleanup package cache

[Timer]
Unit=pkgcaching-server.service
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=timers.target
