# pkgcaching/pkgcaching-server.service
[Unit]
Description=Perform package cache cleanup

[Service]
Type=oneshot
ExecStart=/opt/pkgcache/server.sh
