# btrfs/btrfs-maintain-monthly.service
[Unit]
Description=Run btrfs maintenance scripts

[Service]
Type=oneshot
StandardOutput=null
StandardError=null
ExecStart=/usr/bin/btrfs-maintenance
