# btrfs/btrfs-maintain-monthly.service
[Unit]
Description=Run btrfs maintenance scripts

[Service]
Type=oneshot
ExecStart=/usr/bin/btrfs-maintenance
