# zfs/zfs-scrub-monthly.service
[Unit]
Description=Run zfs scrub

[Service]
Type=oneshot
ExecStart=/usr/bin/zpool scrub tank
